# ops/vps/backups/convy-restore-verify.timer
[Unit]
Description=Run Convy PostgreSQL restore verification weekly

[Timer]
OnCalendar=Sun 04:20:00
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
